# debian/DEBIAN/postinst
#!/bin/sh
update-rc.d -f multipathd-udev-monitor defaults
service multipathd-udev-monitor start
