# debian/DEBIAN/prerm
#!/bin/sh
update-rc.d -f multipathd-udev-monitor remove 
